# packages/peach-ssb-server/debian/peach-ssb-server@.service
[Unit]
Description=A Scuttlebutt server for PeachCloud on %I
Requires=peach-ssb-server.service
Before=peach-ssb-server.service
BindsTo=peach-ssb-server.service

[Service]
Type=simple
EnvironmentFile=/etc/default/peach-ssb-server
Environment=PEACH_CONFIG=/etc/peach/agents/%i/
Environment=PEACH_STATE=/var/lib/peach/agents/%i/
ExecStart=/usr/bin/peach-ssb-server
User=peach
Restart=Always
RestartSec=1
RuntimeMaxSec=3600
StartLimitIntervalSec=0

[Install]
WantedBy=peach-ssb-server.target
